# config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/workspace/logs/supervisord.log
logfile_maxbytes=10MB
pidfile=/tmp/supervisord.pid

[program:agent-server]
command=python3 /workspace/bin/agent-server.py
directory=/workspace
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=45
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:relay]
command=python3 /workspace/bin/relay.py
directory=/workspace
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dashboard]
; Runs pre-built Next.js app. To rebuild after code changes:
; 1. docker exec <container> bash -c "cd /workspace/dashboard && npm run build"
; 2. Then restart: docker restart <container>
command=npx next start -p %(ENV_DASHBOARD_PORT)s
directory=/workspace/dashboard
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:scheduler]
command=python3 /workspace/bin/scheduler.py
directory=/workspace
autostart=true
autorestart=true
stopsignal=TERM
stopwaitsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
